# main/gatus/gatus.initd
#!/sbin/openrc-run

: ${logfile:="/var/log/$RC_SVCNAME"}

command="/usr/bin/gatus"
command_background=yes
command_user=root:root
pidfile="/run/$RC_SVCNAME.pid"
extra_started_commands="reload"
output_log="${logfile}"
error_log="${logfile}"

export GATUS_CONFIG_PATH=${GATUS_CONFIG_PATH:-/etc/gatus.d}

depend() {
	need net
	after firewall
}
